# bootstrap
#!/sbin/openrc-run

description="Bootstrap service"

depend() {
  need localmount
  after bootmisc
}

start() {
  ebegin "Starting bootstrap service(s)"
  /usr/local/bin/bootstrap-startup
  eend 0
}

stop() {
  ebegin "Stopping bootstrap service(s)"
  /usr/local/bin/bootstrap-shutdown
  eend 0
}
